# rc.local
#!/bin/sh -e

bash /root/deploy-agent/start_pxe &
exit 0
