# rc.local
#!/bin/sh -e

bash /root/deploy-agent/start_pxe &
bash /root/deploy-agent/start_agent &
exit 0
