# rc.local
#!/bin/sh -e

bash /home/uniapps/deploy-agent/start_pxe &
exit 0
